# systemd-service/riffbot.service.template
[Unit]
Description=Riffbot - Discord Music Bot

[Service]
WorkingDirectory={PROJECT_ROOT}
ExecStart={VENV_PATH}/bin/python3 -m riffbot
Environment=PYTHONUNBUFFERED=1
Restart=on-failure
{USER_DIRECTIVE}
[Install]
WantedBy=default.target
